# package/scripts/postinstall
#!/bin/zsh
/bin/launchctl bootstrap system /Library/LaunchDaemons/com.munkireport.detectx.hourly_scan.plist
