# package/scripts/postinstall
#!/bin/zsh
/bin/launchctl bootstrap system /LaunchDaemons/org.munkireport.detectx.hourly_scan.plist
